# infra/com.dean-bot.listen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dean-bot.listen</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dean-bot/mac-mini-agent/infra/boot-health-check.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dean-bot/mac-mini-agent/apps/listen</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/dean-bot</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/dean-bot/mac-mini-agent/logs/listen.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dean-bot/mac-mini-agent/logs/listen.err.log</string>
</dict>
</plist>
